# files/etc/systemd/system/cerebro.service
[Unit]
Description=cerebro

[Service]
ExecStart=/opt/cerebro/bin/cerebro

[Install]
WantedBy=multi-user.target
